# test/rc.local
#!/bin/bash -e
initmain.jar
initbash.sh
exit 0
