# test/rc.local
#!/bin/bash -e
/etc/initmain.jar
exit 0
